# delivery/ops/templates/gridflow-app.service
[Unit]
Description=GridFlow Main App
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{APP_USER}}
Group={{APP_GROUP}}
WorkingDirectory={{APP_DIR}}
EnvironmentFile={{ENV_FILE}}
Environment=GRIDFLOW_ENV_FILE={{ENV_FILE}}
ExecStart={{VENV_PYTHON}} -m uvicorn app.main:app --host {{APP_BIND_HOST}} --port {{APP_BIND_PORT}}
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
UMask=0077

[Install]
WantedBy=multi-user.target
