# delivery/ops/templates/gridflow-pg-backup.timer
[Unit]
Description=Run GridFlow PostgreSQL backup daily

[Timer]
OnCalendar=*-*-* 03:20:00 UTC
Persistent=true
RandomizedDelaySec=10m
Unit=gridflow-pg-backup.service

[Install]
WantedBy=timers.target
