# ops/gridflow-health-alert.service
[Unit]
Description=GridFlow basic health alert check
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=root
Group=root
ExecStart=/home/ubuntu/upbit_bot/ops/gridflow_health_alert.sh
SuccessExitStatus=2
